# osquery packs/darwin_survey.conf
{
    "queries": {
        "system_info": {
            "query": "SELECT * from system_info;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Retrieve system information",
            "platform": "darwin"
        },
        "disk_encryption": {
            "query": "select e.*,m.* from disk_encryption as e join mounts as m on e.name=m.device",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Returns information about disk encryption status",
            "platform": "darwin"
        },
        "processes": {
            "query": "SELECT * from processes;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "List process information",
            "platform": "darwin"
        },
        "crontab": {
            "query": "SELECT * from crontab;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Line parsed values from system and user cron/tab",
            "platform": "darwin, linux"
        },
        "launchd": {
            "query": "SELECT * from launchd;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "LaunchAgents and LaunchDaemons from default search paths",
            "platform": "darwin"
        },
        "startup_items": {
            "query": "SELECT * from startup_items;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Applications and binaries set as user/login startup items",
            "platform": "darwin"
        },
        "chrome_extensions": {
            "query": "SELECT * from chrome_extensions;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Chrome extensions",
            "platform": "darwin"
        },
        "certificates": {
            "query": "SELECT * from certificates;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Certificate Authorities installed in Keychains/ca-bundles",
            "platform": "darwin, windows"
        },
        "cups_jobs": {
            "query": "SELECT * from cups_jobs;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Print jobs",
            "platform": "darwin"
        },
        "authorized_keys": {
            "query": "SELECT * from authorized_keys;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "A line-delimited authorized_keys table",
            "platform": "darwin"
        },
        "dns_resolvers": {
            "query": "SELECT * from dns_resolvers;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Resolvers used by this host",
            "platform": "darwin"
        },
        "etc_hosts": {
            "query": "SELECT * from etc_hosts;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Line-parsed /etc/hosts",
            "platform": "darwin"
        },
        "event_taps": {
            "query": "select p.*,e.* from processes as p join event_taps e on e.tapping_process=p.pid;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Returns information about installed event taps",
            "platform": "darwin"
        },
        "firefox_addons": {
            "query": "SELECT * from firefox_addons;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "firefox addons",
            "platform": "darwin"
        },
        "groups": {
            "query": "SELECT * from groups;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "",
            "platform": "darwin"
        },
        "homebrew_packages": {
            "query": "SELECT * from homebrew_packages;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "user groups",
            "platform": "darwin"
        },
        "apps": {
            "query": "SELECT * from apps;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "OS X applications installed in known search paths (e.g., /Applications) ",
            "platform": "darwin"
        },
        "kernel_extensions": {
            "query": "SELECT * from kernel_extensions;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "OS X's kernel extensions, both loaded and within the load search path",
            "platform": "darwin"
        },
        "keychain_items": {
            "query": "SELECT * from keychain_items;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Generic details about keychain items",
            "platform": "darwin"
        },
        "logged_in_users": {
            "query": "SELECT * from logged_in_users;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "logged in users",
            "platform": "darwin"
        },
        "mounts": {
            "query": "SELECT * from mounts;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "mounted file systems",
            "platform": "darwin"
        },
        "nfs_shares": {
            "query": "SELECT * from nfs_shares;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "NFS shares exported by the host",
            "platform": "darwin"
        },
        "preferences": {
            "query": "SELECT * from preferences;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "OS X defaults and managed preferences.",
            "platform": "darwin"
        },
        "python_packages": {
            "query": "SELECT * from python_packages;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Python packages installed in a system",
            "platform": "darwin"
        },
        "routes": {
            "query": "SELECT * from routes;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "The active route table for the host system",
            "platform": "darwin"
        },
        "safari_extensions": {
            "query": "SELECT * from safari_extensions;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Safari browser extension details for all users",
            "platform": "darwin"
        },
        "sandboxes": {
            "query": "SELECT * from sandboxes;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "OS X application sandboxes container details.",
            "platform": "darwin"
        },
        "shared_folders": {
            "query": "SELECT * from shared_folders;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Folders available to others via SMB or AFP",
            "platform": "darwin"
        },
        "sharing_preferences": {
            "query": "SELECT * from sharing_preferences;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "OS X Sharing preferences",
            "platform": "darwin"
        },
        "shell_history": {
            "query": "SELECT * from shell_history;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "A line-delimited (command) table of per-user .*_history data",
            "platform": "darwin"
        },
        "suid_bin": {
            "query": "SELECT * from suid_bin;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "suid binaries in common locations",
            "platform": "darwin"
        },
        "system_extensions": {
            "query": "SELECT * from system_extensions;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "macOS (>= 10.15) system extension table",
            "platform": "darwin"
        },
        "users": {
            "query": "SELECT * from users;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Local user accounts (including domain accounts that have logged on locally (Windows))",
            "platform": "darwin"
        },
        "process_open_files": {
            "query": "SELECT p.*,e.* from processes as p join process_open_files as e USING  (pid);",
            "interval": "86400",
            "version": "1.4.5",
            "description": "File descriptors for each process",
            "platform": "darwin"
        },
        "listening_ports": {
            "query": "SELECT p.*,e.* from processes as p join listening_ports as e USING  (pid);",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Processes with listening (bound) network sockets/ports",
            "platform": "darwin"
        },
        "interface_details": {
            "query": "SELECT a.*,d.* FROM interface_addresses as a join interface_details d on d.interface=a.interface;",
            "interval": "86400",
            "version": "1.4.5",
            "description": "Detailed information and stats of network interfaces",
            "platform": "darwin"
        },
        "alf": {
            "query": "select * from alf;",
            "interval": "86400",
            "platform": "darwin",
            "version": "1.4.5",
            "description": "Retrieves the configuration values for the Application Layer Firewall for OSX.",
            "value": "Verify firewall settings are as expected"
        },
        "alf_exceptions": {
            "query": "select * from alf_exceptions;",
            "interval": "86400",
            "platform": "darwin",
            "version": "1.4.5",
            "description": "Retrieves the exceptions for the Application Layer Firewall in OSX.",
            "value": "Verify firewall settings are as expected"
        },
        "alf_services": {
            "query": "select * from alf_services;",
            "interval": "86400",
            "platform": "darwin",
            "version": "1.4.5",
            "description": "Retrieves the services for the Application Layer Firewall in OSX.",
            "value": "Verify firewall settings are as expected"
        },
        "alf_explicit_auths": {
            "query": "select * from alf_explicit_auths;",
            "interval": "86400",
            "platform": "darwin",
            "version": "1.4.5",
            "description": "Retrieves the list of processes with explicit authorization for the Application Layer Firewall.",
            "value": "Verify firewall settings are as expected"
        }
    }
}